# config/etc/supervisor/conf.d/postfix.conf
[program:postfix]
command=/usr/lib/postfix/master -d
redirect_stderr=true
stopsignal=QUIT
autorestart=true
