# config/etc/supervisor/conf.d/syslogd.conf
[program:syslogd]
command=/sbin/syslogd -n -O /dev/stdout
redirect_stderr=true
stopsignal=QUIT
autorestart=true
